# transcribe-server.service
[Unit]
Description=Granite Speech Transcription API
After=network.target

# Update User and Group in this section
[Service]
Type=simple
User=peter
Group=peter
WorkingDirectory=%h/transcribe-journal
ExecStart=%h/.local/bin/uv run --env-file .env python server/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcribe-server

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
