# transcribe-server.service
[Unit]
Description=Granite Speech Transcription API
After=network.target

# Update User and Group in this section
[Service]
Type=simple
User=peter
Group=peter
WorkingDirectory=/home/peter/transcribe-journal
ExecStart=/home/peter/.local/bin/uv run --env-file .env python server/server.pyRestart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcribe-server

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
