# transcribe-server.service
[Unit]
Description=Granite Speech Transcription API
After=network.target

# Update User and Group in this section
[Service]
Type=simple
User=peter
Group=peter
WorkingDirectory=%h/transcribe-journal
ExecStart=/usr/bin/uv run --env-file .env python server/server.py
Restart=always
RestartSec=10
StandardOutput=append:%h/transcribe-journal/server/logs/server.log
StandardError=append:%h/transcribe-journal/server/logs/error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
