# transcribe-journal.service
[Unit]
Description=Audio Transcription Journal Service
After=network.target

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/coding/transcribe-journal
ExecStart=/usr/bin/uv run --env-file .env watcher.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcribe-journal

[Install]
WantedBy=multi-user.target
